# crontab
* * * * * /usr/local/bin/python /app/cron.py >> /var/log/cron.log 2>&1
